# target/control/deb/postinst
#!/bin/sh

chown -R thingsboard: /var/log/tb-gateway
chown -R thingsboard: /usr/share/tb-gateway
update-rc.d tb-gateway defaults
