# target/control/tb-gateway.service
[Unit]
Description=tb-gateway
After=syslog.target

[Service]
User=thingsboard
ExecStart=/usr/share/tb-gateway/bin/tb-gateway.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
