# target/debian/postrm
#!/bin/sh -e



update-rc.d -f tb-gateway remove


